# etc/systemd/hails-served.service
[Unit]
Description=Regenerate the public requests served counter

[Service]
# Machine specific settings, see config.example.env. The leading dash means a missing file is
# not an error, so the defaults still apply.
EnvironmentFile=-/etc/hails-stats/config.env
Type=oneshot
ExecStart=/usr/local/bin/hails-served.py
