# etc/systemd/hails-verify.timer
[Unit]
Description=Daily verification of the hails events warehouse

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
