# etc/systemd/hails-stats.timer
[Unit]
Description=Regenerate GoAccess dashboard every 15 minutes

[Timer]
OnBootSec=2min
# The warehouse still moves every 5 minutes via hails-ingest.timer. This is only how often the HTML
# is redrawn from it.
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
