# etc/systemd/hails-prune.service
[Unit]
Description=Prune aged events and unreferenced dimension rows from the hails events warehouse

[Service]
EnvironmentFile=-/etc/hails-stats/config.env
Type=oneshot
ExecStart=/usr/local/bin/hails-prune.py
